# systemd/.config/systemd/user/dropbox.service.d/override.conf
[Service]
Type=forking
KillMode=process
PassEnvironment=DISPLAY XAUTHORITY WAYLAND_DISPLAY
ExecStartPre=/usr/bin/sleep 2
